# .dotfiles/crontab/crontab.file
00 4 * * *      $HOME/.dotfiles/crontab/gem_update
10 4 * * *      /usr/bin/env brew update && /usr/bin/env brew install `/usr/bin/env brew outdated --quiet`
