# .dotfiles/crontab/crontab.file
00 4 * * *      $HOME/.dotfiles/crontab/gem_update
10 4 * * *      $HOME/.dotfiles/crontab/brew_update
